# linux/home/.config/systemd/user/borgmatic.timer
[Unit]
Description=Run borgmatic daily

[Timer]
OnActiveSec=1hour
OnUnitActiveSec=1day

