# linux/home/.config/systemd/user/borgmatic.service
[Unit]
Description=Run borgmatic

[Service]
Type=oneshot
ExecStart=/usr/bin/borgmatic

