# linux/home/.config/systemd/user/borgmatic.timer
[Unit]
Description=Run borgmatic daily

[Install]
WantedBy=multi-user.target

[Timer]
OnActiveSec=1hour
OnUnitActiveSec=1day
