# linux/home/.config/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpaper every hour

[Timer]
Persistent=false
OnActiveSec=1sec
OnUnitActiveSec=1hour

