# etc/cron.d/borgmatic
# Run borgmatic nightly at 11PM

0 3 * * * <username> PATH=$PATH:/usr/local/bin /usr/local/bin/borgmatic
