# home/.config/systemd/user/compositor.service
[Unit]
Description=Compositor
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/xcompmgr

[Install]
WantedBy=default.target

