# home/.config/systemd/user/compositor.service
[Unit]
Description=Compositor

[Service]
Type=simple
ExecStart=/usr/bin/picom --experimental-backends

[Install]
WantedBy=xorg.target

