# home/.config/systemd/user/compositor.service
[Unit]
Description=Compositor

[Service]
Type=simple
ExecStart=/usr/bin/xcompmgr

[Install]
WantedBy=xorg.target

