# home/.config/systemd/user/files.service
[Unit]
Description=File watching

[Service]
Type=simple
ExecStart=/home/sam/scripts/cruftwatch.sh

[Install]
WantedBy=default.target

