# home/.config/systemd/user/gamma.service
[Unit]
Description=Gamma management

[Service]
Type=simple
ExecStart=/usr/bin/redshift

[Install]
WantedBy=xorg.target

