# home/.config/systemd/user/gamma.service
[Unit]
Description=Gamma management

[Service]
Type=simple
ExecStart=/usr/bin/redshift -l 51.5941:0.1298

[Install]
WantedBy=xorg.target

