# home/.config/systemd/user/hydroxide.service
[Unit]
Description=Hydroxide

[Service]
Type=simple
ExecStart=/usr/bin/hydroxide serve

[Install]
WantedBy=default.target

