# home/.config/systemd/user/imap.service
[Unit]
Description=IMAP

[Service]
Type=simple
ExecStart=/usr/bin/hydroxide imap

[Install]
WantedBy=default.target

