# home/.config/systemd/user/locker.service
[Unit]
Description=Lock on suspend/sleep

[Service]
Type=simple
ExecStart=/usr/bin/xss-lock -l -- physlock

[Install]
WantedBy=xorg.target

