# home/.config/systemd/user/mpd-mpris.service
[Unit]
Description=mpd-mpris
After=mpd.service

[Service]
Type=simple
ExecStart=/usr/bin/mpd-mpris -no-instance

[Install]
WantedBy=default.target

