# home/.config/systemd/user/notifications.service
[Unit]
Description=Notifications

[Service]
Type=simple
ExecStart=/usr/bin/dunst

[Install]
WantedBy=xorg.target

