# home/.config/systemd/user/search-engine.service
[Unit]
Description=Search engine

[Service]
Type=simple
ExecStart=/usr/bin/bangin-server-node

[Install]
WantedBy=default.target

