# home/.config/systemd/user/smtp.service
[Unit]
Description=SMTP

[Service]
Type=simple
ExecStart=/usr/bin/hydroxide smtp

[Install]
WantedBy=default.target

