# home/.config/systemd/user/wallpaper.service
[Unit]
Description=Re-set persisted wallpaper

[Service]
Type=oneshot
ExecStart=/usr/bin/wal -R

[Install]
WantedBy=xorg.target

