# home/.config/systemd/user/xorg.target
[Unit]
Description=Xorg is up
Requires=default.target
After=default.target

